# sample_deployments/nginx/superverisord.conf
;DefectDojo
[program:defectdojo]
command = /install_path/django-DefectDojo/env/bin/dojo-startup                 ; Command to start app
user = defectdojo                                            ; User to run as
stdout_logfile = /install_path/django-DefectDojo/env/logs/dojo-supervisor.log  ; Where to write log messages
redirect_stderr = true                                 ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8        ; Set UTF-8 as default encoding

;CeleryDefectDojo
[program:celerydefectdojo]
command = /install_path/django-DefectDojo/env/bin/celery-startup                 ; Command to start app
user = defectdojo                                            ; User to run as
stdout_logfile = /install_path/django-DefectDojo/env/logs/celerydojo-supervisor.log  ; Where to write log messages
redirect_stderr = true                                 ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8        ; Set UTF-8 as default encoding
